# .kdm/profile.d/00-base
# bash-env variable set
export BASH_ENV_DIR_RC_D="${BASH_ENV_DIR_KDM}/rc.d"
export BASH_ENV_DIR_CODE="${HOME}/code"
export BASH_ENV_DIR_KDM="${HOME}/.kdm"
export BASH_ENV_DIR_NMAP="${HOME}/.nmap"
export BASH_ENV_DIR_SNMP="${HOME}/.snmp"
export BASH_ENV_DIR_SSH="${HOME}/.ssh"
export BASH_ENV_DIR_SSH_CONFIG="${HOME}/.ssh/config.d"
export BASH_ENV_DIR_SSH_SESSIONS="${HOME}/.ssh-sessions"
export BASH_ENV_DIR_VIM="${HOME}/.vim"
export BASH_ENV_DIR_VIM_BACKUP="${BASH_ENV_DIR_VIM}/backup"
export BASH_ENV_DIR_VIM_TMP="${BASH_ENV_DIR_VIM}/tmp"

export BASH_ENV_DIR_EXTRA="${BASH_ENV_DIR_KDM}/extra"
export BASH_ENV_DIR_EXTRA_LOGOUT="${BASH_ENV_DIR_KDM}/extra-logout"
export BASH_ENV_DIR_PROFILE_D="${BASH_ENV_DIR_KDM}/profile.d"
export BASH_ENV_FILE_HISTORY="${BASH_ENV_DIR_KDM}/history"

export BASH_ENV_FILE_LESSFILTER="${HOME}/.lessfilter"
export BASH_ENV_FILE_LOGOUT="${HOME}/.bash_logout"
export BASH_ENV_FILE_PROFILE="${BASH_ENV_DIR_KDM}/profile"
export BASH_ENV_FILE_RC="${BASH_ENV_DIR_KDM}/rc"
export BASH_ENV_FILE_COLORS="${BASH_ENV_DIR_PROFILE_D}/03-colors"
export BASH_ENV_FILE_CONFIG="${BASH_ENV_DIR_KDM}/bash-env.conf"
export BASH_ENV_FILE_DIRCOLORS="${HOME}/.dircolors"
export BASH_ENV_FILE_GEMRC="${HOME}/.gemrc"
export BASH_ENV_FILE_GITCONFIG="${HOME}/.gitconfig"
export BASH_ENV_FILE_GITIGNORE="${HOME}/.gitignore"
export BASH_ENV_FILE_INPUTRC="${HOME}/.inputrc"
export BASH_ENV_FILE_NLOAD="${HOME}/.nload"
export BASH_ENV_FILE_PIPRC="${HOME}/.piprc"
export BASH_ENV_FILE_SCREENRC="${HOME}/.screenrc"
export BASH_ENV_FILE_SSH_AUTHKEYS="${BASH_ENV_DIR_SSH}/authorized_keys"
export BASH_ENV_FILE_SSH_CACHE="${BASH_ENV_DIR_KDM}/sshcache"
export BASH_ENV_FILE_SSH_CONFIG="${BASH_ENV_DIR_SSH}/config"
export BASH_ENV_FILE_SUDOCACHE="${BASH_ENV_DIR_KDM}/sudocache"
export BASH_ENV_FILE_TIDY="${HOME}/.tidyrc"
export BASH_ENV_FILE_TOPRC="${HOME}/.toprc"
export BASH_ENV_FILE_VIMRC="${HOME}/.vimrc"

# Custom pip config file path
export PIP_CONFIG_FILE="${BASH_ENV_FILE_PIPRC}"

# Custom tidy-html5 config file path
export HTML_TIDY="${BASH_ENV_FILE_TIDY}"


# vim: set syntax=sh filetype=sh ts=2 sw=2 tw=0 noet :
