# .kdm/profile.d/00-base
# bash-env variable set
export BASH_ENV_DIR_RC_D="${BASH_ENV_DIR_KDM}/rc.d"
export BASH_ENV_DIR_CODE="${HOME}/code"
export BASH_ENV_DIR_KDM="${HOME}/.kdm"
export BASH_ENV_DIR_NMAP="${HOME}/.nmap"
export BASH_ENV_DIR_SNMP="${HOME}/.snmp"

export BASH_ENV_DIR_SSH="${HOME}/.ssh"
export BASH_ENV_DIR_SSH_CONFIG="${HOME}/.ssh/config.d"
export BASH_ENV_DIR_SSH_SESSIONS="${HOME}/.ssh-sessions"

export BASH_ENV_DIR_EXTRA="${BASH_ENV_DIR_KDM}/extra"
export BASH_ENV_DIR_EXTRA_LOGOUT="${BASH_ENV_DIR_KDM}/extra-logout"
export BASH_ENV_DIR_PROFILE_D="${BASH_ENV_DIR_KDM}/profile.d"

export BASH_ENV_FILE_HISTORY="${BASH_ENV_DIR_KDM}/history"


export BASH_ENV_FILE_LESSFILTER="${HOME}/.lessfilter"
export BASH_ENV_FILE_LOGOUT="${HOME}/.bash_logout"
export BASH_ENV_FILE_PROFILE="${BASH_ENV_DIR_KDM}/profile"
export BASH_ENV_FILE_RC="${BASH_ENV_DIR_KDM}/rc"
export BASH_ENV_FILE_COLORS="${BASH_ENV_DIR_PROFILE_D}/03-colors"
export BASH_ENV_FILE_CONFIG="${BASH_ENV_DIR_KDM}/bash-env.conf"
export BASH_ENV_FILE_DISABLE_TRUECOLOR="${BASH_ENV_DIR_KDM}/disable-truecolor"
export BASH_ENV_FILE_DIRCOLORS="${HOME}/.dircolors"
export BASH_ENV_FILE_GEMRC="${HOME}/.gemrc"
export BASH_ENV_FILE_GITCONFIG="${HOME}/.gitconfig"
export BASH_ENV_FILE_GITIGNORE="${HOME}/.gitignore"
export BASH_ENV_FILE_INPUTRC="${HOME}/.inputrc"
export BASH_ENV_FILE_NLOAD="${HOME}/.nload"
export BASH_ENV_FILE_PIPRC="${HOME}/.piprc"
export BASH_ENV_FILE_SCREENRC="${HOME}/.screenrc"
export BASH_ENV_FILE_SSH_AUTHKEYS="${BASH_ENV_DIR_SSH}/authorized_keys"
export BASH_ENV_FILE_SSH_CACHE="${BASH_ENV_DIR_KDM}/sshcache"
export BASH_ENV_FILE_SSH_CONFIG="${BASH_ENV_DIR_SSH}/config"
export BASH_ENV_FILE_SUDOCACHE="${BASH_ENV_DIR_KDM}/sudocache"
export BASH_ENV_FILE_TIDY="${HOME}/.tidyrc"
export BASH_ENV_FILE_TOPRC="${HOME}/.toprc"
export BASH_ENV_FILE_VIMRC="${HOME}/.vimrc"


# Custom pip config file path
export PIP_CONFIG_FILE="${BASH_ENV_FILE_PIPRC}"

# Custom tidy-html5 config file path
export HTML_TIDY="${BASH_ENV_FILE_TIDY}"


# Try to load config from facter
_truecolor_config() {
	[[ "${TERM}" == "xterm-256color+iterm3" ]] && COLORTERM="truecolor"

	[[ "${COLORTERM}" == "24bit" ]] && COLORTERM="truecolor"

	[[ -e "${BASH_ENV_FILE_DISABLE_TRUECOLOR}" ]] && unset COLORTERM

	# Extrapolate ${COLORTERM} into boolean ${COLOR_16M}
	case "${COLORTERM}" in
		"truecolor") COLOR_16M="true"  ;;
		*)           COLOR_16M="false" ;;
	esac

	# if [[ "${COLOR_16M}" == "true" ]]; then
	# 	# Physical console on Linux:
	# 	# Disable true color (where as it would be enabled if connecting via SSH or the like)
	# 	if tty 2> /dev/null | grep -Eq '^\/dev\/tty[0-9]$'; then
	# 		COLOR_16M="false"
	# 	fi
	# fi

	export COLOR_16M
}

_prompt_config() {
	if ! hash facter > /dev/null 2>&1; then
		PROMPT_HOSTNAME="true"
		PROMPT_USERNAME="true"
		PROMPT_TIMEDATE="true"

		export PROMPT_HOSTNAME PROMPT_USERNAME PROMPT_TIMEDATE
		return
	fi

	# PROMPT_HOSTNAME="$(facter bash-env.prompt.hostname 2> /dev/null)"
	# PROMPT_USERNAME="$(facter bash-env.prompt.username 2> /dev/null)"
	# PROMPT_TIMEDATE="$(facter bash-env.prompt.timedate 2> /dev/null)"

	[[ "${PROMPT_HOSTNAME}" != "false" ]] && PROMPT_HOSTNAME="true"
	[[ "${PROMPT_USERNAME}" != "false" ]] && PROMPT_USERNAME="true"
	[[ "${PROMPT_TIMEDATE}" != "false" ]] && PROMPT_TIMEDATE="true"

	# Override to always show hostname when connected via SSH
	[[ -z "${SSH_CLIENT+x}"     ]] || PROMPT_HOSTNAME="true"
	[[ -z "${SSH_CONNECTION+x}" ]] || PROMPT_HOSTNAME="true"
	[[ -z "${SSH_TTY+x}"        ]] || PROMPT_HOSTNAME="true"

	export PROMPT_HOSTNAME PROMPT_USERNAME PROMPT_TIMEDATE
}


_truecolor_config
_prompt_config


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
