# .kdm/profile.d/02-arrays
# shellcheck disable=SC2034

# kdm bash-env bash_profile arrays

# Ignore bland commands
ARRAY_HISTIGNORE=(
..
[bf]g
c
cd
clear
cs
dc
dl
dt
exit
g-gca
g-gcp
g-gi
g-gs
g-gu
go2sleep
l
l.
l1
la
ll
ll.
lld
lll
logout
ls
stayup
)
# Generate HISTIGNORE from array
STR_HISTIGNORE="${ARRAY_HISTIGNORE[*]}"
HISTIGNORE="${STR_HISTIGNORE//\ /:}"


# Array of editor aliases
ARRAY_EDITOR_ALIASES=(v vi vim nvim)

# Array of editor binaries
ARRAY_EDITOR_NAMES=(nvim vim vi)


# Detect if /usr/local (Intel) or /opt/homebrew (Apple Silicon) should be used as the prefix
PREFIX_ARM="/opt/homebrew"
PREFIX_X86="/usr/local"

PREFIX="${PREFIX_X86}"
if [[ "$(sysctl -n hw.optional.arm64 2> /dev/null)" == "1" ]]; then
	PREFIX="${PREFIX_ARM}"
fi

# Dir shortcuts
UL="${PREFIX}"
ULO="${UL}/opt"
ULR="${UL}/rvm"


# Entries at the beginning of this array will end up at the beginning of MANPATH
# Entries at the end       of this array will end up at the end       of MANPATH
ARRAY_MANPATH=(
"${ULO}/coreutils/libexec/man"
"${ULO}/curl/share/man"
"${ULO}/diffutils/share/man"
"${ULO}/findutils/libexec/man"
"${ULO}/gawk/share/man"
"${ULO}/gcc/share/man"
"${ULO}/gnu-sed/libexec/man"
"${ULO}/gnu-tar/libexec/man"
"${ULO}/gnu-which/libexec/man"
"${ULO}/grep/libexec/man"
"${ULO}/libtool/libexec/man"
"${ULO}/make/libexec/man"
"${ULO}/openssh/share/man"
"${ULO}/openssl@3/share/man"
"${ULO}/openssl@1.1/share/man"
"${ULO}/openssl/share/man"
"${UL}/share/man"
/usr/pgsql-12/share/man
/usr/share/man/overrides
/usr/share/man/en
/usr/share/man
/usr/X11R6/share/man
/opt/X11/share/man
)

# Entries at the beginning of this array will end up at the beginning of MIBDIRS
# Entries at the end       of this array will end up at the end       of MIBDIRS
ARRAY_MIBDIRS=(
"${BASH_ENV_DIR_SNMP}/mibs"
/usr/share/snmp/mibs
)

# Entries at the beginning of this array will end up at the beginning of NODE_PATH
# Entries at the end       of this array will end up at the end       of NODE_PATH
ARRAY_NODE_PATH=(
"${ULO}/node/lib/node_modules"
"${UL}/lib/node_modules"
/usr/lib/node_modules
/usr/lib/nodejs
)

# Entries at the beginning of this array will end up at the beginning of PATH
# Entries at the end       of this array will end up at the end       of PATH
mapfile -t ARRAY_PATH < <(echo -e "${PATH//\:/\\n}")

ARRAY_PATH=(
"${ARRAY_PATH[@]}"
"${HOME}/bin"
"${HOME}/.local/bin"
"${HOME}/.go/bin"
"${HOME}/.composer/vendor/bin"
"${HOME}/perl5/bin"
"${HOME}/Library/Android/sdk/platform-tools"
"${HOME}/.jenv/shims"
'/data/data/com.termux/files/usr/bin/applets'
'/data/data/com.termux/files/usr/bin'
"${ULR}/gems/ruby-2.4.0/bin"
"${ULR}/gems/ruby-2.4.0@global/bin"
"${ULR}/rubies/ruby-2.4.0/bin"
"${ULR}/bin"
"${ULO}/perl/bin"
"${ULO}/net-snmp/bin"
"${ULO}/go/libexec/bin"
"${ULO}/coreutils/libexec/gnubin"
"${ULO}/ed/libexec/gnubin"
"${ULO}/findutils/libexec/gnubin"
"${ULO}/gawk/libexec/gnubin"
"${ULO}/gnu-indent/libexec/gnubin"
"${ULO}/gnu-sed/libexec/gnubin"
"${ULO}/gnu-tar/libexec/gnubin"
"${ULO}/grep/libexec/gnubin"
"${ULO}/make/libexec/gnubin"
"${ULO}/curl-openssl/bin"
"${ULO}/curl/bin"
"${ULO}/gnu-getopt/bin"
"${ULO}/gnu-which/bin"
"${ULO}/ncurses/bin"
"${ULO}/python@3.10/bin"
"${ULO}/python@3.10/libexec/bin"
"${ULO}/python@3.10/Frameworks/Python.framework/Versions/3.10/bin"
"${ULO}/python@3.9/bin"
"${ULO}/python@3.9/libexec/bin"
"${ULO}/python@3.9/Frameworks/Python.framework/Versions/3.9/bin"
"${ULO}/python@3.8/bin"
"${ULO}/python@3.8/libexec/bin"
"${ULO}/python@3.8/Frameworks/Python.framework/Versions/3.8/bin"
"${ULO}/python@3.7/bin"
"${ULO}/python@3.7/libexec/bin"
"${ULO}/python@2/bin"
"${ULO}/python/libexec/bin"
"${ULO}/libxml2/bin"
"${ULO}/libxslt/bin"
"${ULO}/openssl@3/bin"
"${ULO}/openssl@1.1/bin"
"${ULO}/openssl/bin"
"${ULO}/ruby/bin"
'/opt/local/sbin'
'/opt/local/bin'
'/opt/puppetlabs/bin'
"${UL}/sbin"
"${UL}/bin"
'/System/Library/Filesystems/apfs.fs/Contents/Resources'
'/usr/pgsql-12/bin'
'/usr/pgsql-11/bin'
'/usr/pgsql-10/bin'
'/System/Library/CoreServices'
'/Library/Apple/usr/libexec'
'/Library/Apple/usr/bin'
'/usr/libexec'
'/usr/bin/X11'
'/usr/bin/site_perl'
'/usr/bin/vendor_perl'
'/usr/bin/core_perl'
'/opt/X11/bin'
'/usr/X11R6/bin'
'/usr/games'
'/usr/lib/mit/bin'
'/usr/lib/mit/sbin'
'/usr/bin'
'/bin'
'/usr/sbin'
'/sbin'
'/opt/cuda/bin'
)
# '/Applications/VMware Fusion.app/Contents/Library'
# '/Applications/VMware Fusion.app/Contents/Public'


ARRAY_MKDIR=(
"${BASH_ENV_DIR_KDM}"
"${BASH_ENV_DIR_SSH_CONFIG}"
"${BASH_ENV_DIR_SSH_SESSIONS}"
)

# These are intentionally not ordered alphabetically
ARRAY_PACKAGE_MANAGERS=(
brew
pacman
apt
apt-get
dnf
dnf-3
yum
zypper
)

# Array of shell options to be set
ARRAY_SHELL_OPTIONS=(
cdable_vars
cdspell
checkhash
checkwinsize
histappend
no_empty_cmd_completion
)
# Removed for multicolor PS1
# xpg_echo

# Array of shell options that only work with bash 4+
ARRAY_SHELL_OPTIONS_BASH4=(
autocd
dirspell
)

# Sourcing entire bash_completion library is slow,
# about 0.5s on i7-9700k@5.2GHz, 3200MHz DDR4, and NVME flash
ARRAY_SOURCE=(
# /etc/profile.d/bash_completion.sh
# "${UL}/etc/profile.d/bash_completion.sh"
"${BASH_ENV_FILE_RC}"
)

# Array of automatic sudo-ing commands (Linux and macOS)
ARRAY_SUDO=(
arp
arp-scan
dd
dhcping
dmesg
iftop
ioping
jnettop
mount
netstat
nload
noping
oping
reboot
route
shutdown
sntop
su
sysctl
syslog
umount
visudo
)

# Array of automatic sudo-ing commands (Linux-only)
ARRAY_SUDO_LINUX=(
apt
apt-cache
apt-get
dnf
dnf-3
ether-wake
etherwake
fdisk
fedup
journalctl
pacman
rpi-update
rpm
service
systemctl
virsh
virt-admin
virt-clone
virt-convert
virt-host-validate
virt-install
virt-manager
virt-pki-validate
virt-viewer
virt-xml
virt-xml-validate
yum
zypper
)

# Array of automatic sudo-ing commands (macOS-only)
ARRAY_SUDO_MACOS=(
apfs.util
atsutil
chflags
cupsctl
defaults
diskutil
dscacheutil
dspci
fs_usage
gpt
installer
ioreg
iostat
kextcache
kextload
kextstat
kextunload
kextutil
keytool
kmutil
launchctl
log
mdutil
newfs_apfs
newfs_exfat
newfs_hfs
newfs_msdos
newfs_udf
nvram
pfctl
pmset
powermetrics
productutil
puppet
reboot
repair_packages
rtcread
scutil
security
shutdown
smc
smcread
softwareupdate
spctl
ssdtprgen
sysctl
system_profiler
systemsetup
systemstats
tmutil
trimforce
xcode-select
xcodebuild
)

ARRAY_TOUCH=(
"${BASH_ENV_FILE_HISTORY}"
"${BASH_ENV_FILE_SSH_AUTHKEYS}"
"${BASH_ENV_FILE_SSH_CACHE}"
"${BASH_ENV_FILE_SSH_CONFIG}"
"${BASH_ENV_FILE_SUDOCACHE}"
"${BASH_ENV_DIR_SSH_CONFIG}/base.conf"
)


# Export above array variables
export ARRAY_SUDO
export ARRAY_SUDO_LINUX
export ARRAY_SUDO_MACOS
export HISTIGNORE


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
