# .kdm/profile.d/09-less
# less environment variables to enable color escape sequences and use ${HOME}/.lessfilter
LESS="--raw-control-chars --tilde --shift=10 --quit-on-intr --squeeze-blank-lines --tabs=2"
# --quit-if-one-screen

LESSOPEN="|${BASH_ENV_FILE_LESSFILTER} %s"
export LESS LESSOPEN


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
