# .kdm/profile.d/11-uname
# shellcheck disable=SC1090

# Machine uname variables
if [[ "${BASH_ENV_CACHE_UNAME}" != "true" ]]; then
	cat <<- EOF > "${BASH_ENV_FILE_CACHE_UNAME}"
	export UNAME_KERNEL_NAME="$(uname -s)"
	export UNAME_KERNEL_RELEASE="$(uname -r)"
	export UNAME_MACHINE="$(uname -m)"
	export BASH_ENV_CACHE_UNAME=true
	EOF
fi

. "${BASH_ENV_FILE_CACHE_UNAME}"

# Commented these out as they are unused (for now)
# UNAME="$(               uname)"
# UNAME_ALL="$(           uname -a)"
# UNAME_KERNEL_VERSION="$(uname -v)"
# UNAME_NODE="$(          uname -n)"
# UNAME_PROCESSOR="$(     uname -p)"
# export UNAME UNAME_ALL UNAME_KERNEL_VERSION UNAME_MACHINE UNAME_NODE UNAME_PROCESSOR


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
