# .kdm/profile.d/11-uname
# Machine uname variables
UNAME_KERNEL_NAME="$(   uname -s)"
UNAME_KERNEL_RELEASE="$(uname -r)"
UNAME_MACHINE="$(       uname -m)"

export UNAME_KERNEL_NAME UNAME_KERNEL_RELEASE UNAME_MACHINE

# Commented these out as they are unused (for now)
# UNAME="$(               uname)"
# UNAME_ALL="$(           uname -a)"
# UNAME_KERNEL_VERSION="$(uname -v)"
# UNAME_NODE="$(          uname -n)"
# UNAME_PROCESSOR="$(     uname -p)"
# export UNAME UNAME_ALL UNAME_KERNEL_VERSION UNAME_MACHINE UNAME_NODE UNAME_PROCESSOR


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
