# .kdm/profile.d/13-tty
CURRENT_TTY="$(ps ax | awk '/'"${$}"'/ {printf $2; exit}')"
GPG_TTY="$(tty)"

export CURRENT_TTY GPG_TTY

# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
