# .kdm/profile.d/13-tty
THIS_TTY="$(tty)"

CURRENT_TTY="${THIS_TTY//\/dev\/tty}"
GPG_TTY="${THIS_TTY}"

unset THIS_TTY
export CURRENT_TTY GPG_TTY


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
