# .kdm/profile.d/14-pkgman
# Detect package manager
for PACKAGE_MANAGER in "${ARRAY_PACKAGE_MANAGERS[@]}"; do
	hash "${PACKAGE_MANAGER}" && break
	unset PACKAGE_MANAGER
done
[[ -n "${PACKAGE_MANAGER}" ]] && export PACKAGE_MANAGER


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
