# .kdm/profile.d/15-jenv
# jenv setup to manage multiple java versions
# first, check if jenv is installed
if hash jenv; then
	# Java home variable on OS X
	[[ -d /usr/local/var/jenv ]] && export JENV_ROOT="/usr/local/var/jenv"

	# Load jenv automatically
	eval "$(jenv init - --no-rehash)"

	[[ "${BASH_SOURCE[0]}" != "${0}" ]] && return 0 || exit 0
fi

# If no jenv, setup JAVA_HOME variable.
# Check if the java_home binary is present
if [[ -s /usr/libexec/java_home ]]; then
	# Send all to stdout and grep for filepath
	JAVA_HOME="$(/usr/libexec/java_home 2>&1 | grep '\/')"

	# If the return is a path (aka not empty) - export it; otherwise, unset it
	if [[ "${JAVA_HOME}" ]]; then
		export JAVA_HOME
	else
		unset JAVA_HOME
	fi

	[[ "${BASH_SOURCE[0]}" != "${0}" ]] && return 0 || exit 0
fi


# vim: set syntax=sh filetype=sh ts=2 sw=2 tw=0 noet :
