# .kdm/profile.d/20-init
# Create directories if missing
for ENTRY in "${ARRAY_MKDIR[@]}"; do
	[[ ! -d "${ENTRY}" ]] && mkdir -p "${ENTRY}"
done

# Touch files
for ENTRY in "${ARRAY_TOUCH[@]}"; do
	[[ ! -e "${ENTRY}" ]] && touch "${ENTRY}"
done

echo -n

# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
