# .kdm/profile.d/20-prompt
# shellcheck disable=SC2016

# _prompt_exitcode() outputs failed exit codes if cmd exits non-zero
# Trap current and previous commands to access them

# Check exitcode of previous command
# Always write lines to history file on each prompt
# Generate / export PS1 prompt
# https://unix.stackexchange.com/a/18443

_prompt_command_capture() {
	export CMD_CRNT="${CMD_PRVS}"
	export CMD_PRVS="${1}"
}


preexec_functions+=('_prompt_command_capture')

precmd_functions+=('_prompt_exitcode "${__bp_last_ret_value}"')
precmd_functions+=('history -a')
precmd_functions+=('_prompt_generate')

trap 'CMD_PRVS="${CMD_CRNT}"; CMD_CRNT="${BASH_COMMAND}"' DEBUG

PROMPT_COMMAND+=('_prompt_exitcode "${?}"')
PROMPT_COMMAND+=('history -a')
PROMPT_COMMAND+=('_prompt_generate')


# Trim \w in PS1 to 2 directories (Bash 4+)
export PROMPT_DIRTRIM="2"


# vim: set filetype=sh ts=2 sw=2 tw=0 noet :
